# packages/mediacenter/remotepi-board/system.d/remotepi-board.service
[Unit]
Description=RemotePi Board service
ConditionPathExists=/storage/.cache/services/remotepi-board.conf
Before=plexht.service
After=graphical.target

[Service]
Type=simple
Environment=HOME=/storage
ExecStart=/usr/bin/remotepi-board

[Install]
WantedBy=plexht.service