# packages/mediacenter/remotepi-board/system.d/remotepi-board-shutdown.service
[Unit]
Description=RemotePi Board shutdown service
ConditionPathExists=/storage/.cache/services/remotepi-board.conf
After=plexht.service
Before=systemd-poweroff.service
DefaultDependencies=no

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=/usr/bin/remotepi-board-shutdown

[Install]
WantedBy=poweroff.target